# nota-bot.service
[Unit]
Description=Nota Telegram Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Nota-v3
Environment="PYTHONPATH=/root/Nota-v3"
ExecStartPre=/bin/bash -c 'pkill -f "python.*bot\.py" || true'
ExecStartPre=/bin/rm -f /root/Nota-v3/bot.pid
ExecStart=/root/Nota-v3/venv/bin/python3 /root/Nota-v3/bot.py --force-restart
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3
StandardOutput=append:/root/Nota-v3/logs/bot_stdout.log
StandardError=append:/root/Nota-v3/logs/bot_stderr.log

[Install]
WantedBy=multi-user.target 